# libs/vedana-backoffice/deploy/single-port-container/supervisord.conf
[supervisord]
nodaemon=true

[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true

[program:backend]
command=reflex run --env prod --backend-only
directory=/app/projects/vedana
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
