# libs/vedana-backoffice/deploy/single-port-container/services.d/caddy/run
#!/command/with-contenv bash
exec caddy run --config /etc/caddy/Caddyfile