# libs/vedana-backoffice/deploy/single-port-container/services.d/reflex-backend/run
#!/command/with-contenv bash
cd /app/projects/vedana
exec reflex run --env prod --backend-only